# docker/cruncher/script-publisher.conf
[program:script-publisher]
command=php /var/www/bin/script_publisher.php worker:start script_publisher 1
startsecs=0
autorestart=1
redirect_stderr=true
